# agent-runners/quadlets/langflow.container
[Unit]
Description=langflow container
Requires=lab-network.service
Requires=postgres.service
Requires=traefik.service
After=postgres.service

[Container]
ContainerName=langflow
Image=docker.io/langflowai/langflow:latest
Network=lab.network
Volume=langflow.volume:/app/langflow:Z
Environment=LANGFLOW_DATABASE_URL=postgresql://langflow:password@postgres:5432/langflow
Environment=LANGFLOW_CONFIG_DIR=/app/langflow
Label=traefik.enable=true
Label=traefik.http.routers.langflow.rule=Host(`langflow.vx`)
Label=traefik.http.routers.langflow.entrypoints=websecure
Label=traefik.http.routers.langflow.tls=true
Label=traefik.http.services.langflow.loadbalancer.server.port=7860

[Install]
WantedBy=multi-user.target